# scripts/com.ollie.options.maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Headless cohort-filling maintenance for the options screener.

  Template: __REPO__ and __PYTHON__ are substituted by scripts/install_launchagent.sh.
  Fires inside each of the three auto-log windows (see src/maintenance.autolog_window):
  ds 10:15-11:30, sps 12:15-13:30, ics 14:00-15:00, Mon-Fri. The maintenance
  throttle makes repeat fires idempotent (once per window per day).

  IMPORTANT (macOS): after installation the job will report status 78 (EX_CONFIG)
  until you approve it under System Settings -> General -> Login Items & Extensions
  -> Allow in the Background.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ollie.options.maintenance</string>
    <key>WorkingDirectory</key>
    <string>__REPO__</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>-m</string>
        <string>src.maintenance</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>__REPO__</string>
    </dict>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>20</integer><key>Weekday</key><integer>1</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>20</integer><key>Weekday</key><integer>2</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>20</integer><key>Weekday</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>20</integer><key>Weekday</key><integer>4</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>20</integer><key>Weekday</key><integer>5</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>20</integer><key>Weekday</key><integer>1</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>20</integer><key>Weekday</key><integer>2</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>20</integer><key>Weekday</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>20</integer><key>Weekday</key><integer>4</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>20</integer><key>Weekday</key><integer>5</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>5</integer><key>Weekday</key><integer>1</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>5</integer><key>Weekday</key><integer>2</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>5</integer><key>Weekday</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>5</integer><key>Weekday</key><integer>4</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>5</integer><key>Weekday</key><integer>5</integer></dict>
    </array>
    <key>StandardOutPath</key>
    <string>__REPO__/logs/launchagent.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/logs/launchagent.log</string>
</dict>
</plist>
